# system/crontab
0   *	  *	  *	  *   /home/lmullen/acad/dissertation/wordcount.sh
# *   */2 *   *   *   /bin/sh /Users/lmullen/dev/converts-control/download-omeka-backup.sh
